# crontab
* * * * * cd /var/lib/openshift/544797c25973cab88a000d14/app-deployments/current/repo; ./manage.py updatefeed
