# filescan/filescan.cron
* * * * * cd /opt/currencyfeed; python -c 'from filescan import tasks; tasks.scantask()'
